# runit/rc.local
#                           _           ___  __  __ ____  
#                          / \   _ __  / _ \|  \/  |  _ \ 
#                         / _ \ | '_ \| | | | |\/| | |_) |
#                        / ___ \| | | | |_| | |  | |  _ < 
#                       /_/   \_\_| |_|\___/|_|  |_|_| \_\
if [ -n "$(acpi -a|rg on)" ]
then
	AnOMR switch hybrid
elif [ -n "$(acpi -a|rg off)" ]
then
	AnOMR switch intel
fi
